# userpatches/overlay/auto_switch_usb/usb-role-switch.timer
[Unit]
Description=Timer that checks for usb direction changes                       

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
