# userpatches/overlay/auto_switch_usb/usb-role-switch.service
[Unit]
Description=Switches USB role

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-role-switch
LogLevelMax=warning
