# userpatches/overlay/auto_disable_ssh/auto-disable-ssh.service
[Unit]
Description=Disables SSH if password has not been changed

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-disable-ssh
